# src/docker/crontab/mirror
0 */4 * * * /mirror.sh >> /var/log/cron.log 2>&1
